# macos/com.xerotier.xim-agent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xerotier.xim-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>@ENTRYPOINT@</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>@HOME@</string>
        <key>PATH</key>
        <string>@VENV_BIN@:@HOME@/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>XEROTIER_AGENT_JOIN_KEY</key>
        <string>@JOIN_KEY@</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>StandardOutPath</key>
    <string>@HOME@/Library/Logs/xerotier/xim-agent.out.log</string>
    <key>StandardErrorPath</key>
    <string>@HOME@/Library/Logs/xerotier/xim-agent.err.log</string>
</dict>
</plist>
